# extras/etc/cron.daily/run_md5_array
#!/usr/bin/bash

PATH="$PATH:$ALEXBIN" $ALEXBIN/touch_secrets
PATH="$PATH:$ALEXBIN" $ALEXBIN/md5recurse_unraid
