# extras/etc/cron.daily/run_md5_array
#!/usr/bin/bash

PATH="$PATH:/boot/bin" /boot/bin/touch_secrets
PATH="$PATH:/boot/bin" /boot/bin/md5recurse_unraid
